# samples/yeelight-redshift.timer
[Unit]
Description=Yeelight x Redshift

[Timer]
OnUnitActiveSec=1200s
OnBootSec=1200s

[Install]
WantedBy=multi-user.target 
